# survey/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:survey]
directory=/data/www/survey
command=/usr/bin/python2.7 manage.py run_gunicorn -c /etc/default/gunicorn --settings=survey.settings.production

