# survey/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:survey]
directory=/data/www/survey
#command=gunicorn -k gevent -b 0.0.0.0:8000 -w 2 --max-requests 4096 application:app

