# docker/celery/supervisord.conf
[supervisord]
nodaemon=true

[program:nut]
directory=/data/www/nut
command=/usr/bin/python2.7 /data/www/nut/manage.py celeryd -E -c 8 --time-limit=300 -f /tmp/celeryd.log -l ERROR
