# spider/engine/supervisord.conf
[supervisord]
nodaemon=true

[program:rambutan]
directory=/data/www/spider
command=scrapyd --logfile=/tmp/spider.log --pidfile=spider.pid

