# dew/engine/supervisord.conf
[supervisord]
nodaemon=true

[program:dew]
directory=/data/www/dew
command=gunicorn -k gevent -b 0.0.0.0:5000 -w 4 --max-requests 1024 applications:app