# mangosteen/engine/supervisord.conf
[supervisord]
nodaemon=true

[program:mangosteen]
directory=/data/www/mangosteen
command=gunicorn -k gevent -b 0.0.0.0:5000 -w 4 --max-requests 4096 applications:app

